# desktop/systemd/dms-greeter-resync.service
[Unit]
Description=Sincroniza o wallpaper/tema do greeter do DMS

[Service]
Type=oneshot
ExecStart=%h/.config/DankMaterialShell/greeter-resync.sh
